# lxd-nslookup.service
[Unit]
Description=LXD DNS Lookup Service
After=network.target network-online.target snap.lxd.daemon.service
Wants=network-online.target
Requires=snap.lxd.daemon.service

[Service]
Type=exec
User=lxd
Group=lxd
WorkingDirectory=/run/lxd_nslookup
ExecStart=/usr/local/bin/lxd_nslookupd -c /etc/lxd_nslookup.yml
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StartLimitInterval=0

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
LockPersonality=true

# Network restrictions
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
IPAccounting=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lxd-nslookup

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1

[Install]
WantedBy=multi-user.target
